# edge/crontab
*/5 * * * * www-data /var/pmc/update_meta.sh
