# edge/pmc-restart.service
# Restart the PMC edge nginx server if it's failed.
[Unit]
Description=PMC restart

[Service]
Type=oneshot
ExecStart=/var/pmc/restart-nginx.sh
